# roles/ceph-mon/templates/ceph-mon.service.j2
[Unit]
Description=Ceph Monitor
After=docker.service

[Service]
EnvironmentFile=-/etc/environment
ExecStartPre=-/usr/bin/docker rm ceph-mon-%i
ExecStartPre=$(command -v mkdir) -p /etc/ceph /var/lib/ceph/mon
ExecStart=/usr/bin/docker run --rm --name ceph-mon-%i --net=host \
   {% if not containerized_deployment_with_kv -%}
   -v /var/lib/ceph:/var/lib/ceph \
   -v /etc/ceph:/etc/ceph \
   {% else -%}
   -e KV_TYPE={{kv_type}} \
   -e KV_IP={{kv_endpoint}}\
   -e KV_PORT={{kv_port}} \
   {% endif -%}
   -v /etc/localtime:/etc/localtime:ro \
   {% if mon_docker_privileged -%}
   --privileged \
   {% endif -%}
   {% if mon_docker_net_host -%}
   --net=host \
   {% endif -%}
   -e CEPH_DAEMON=MON \
   -e CEPH_PUBLIC_NETWORK={{ ceph_mon_docker_subnet }} \
   -e IP_VERSION={{ ip_version[-1:] }} \
   {% if monitor_address is defined and monitor_address != '0.0.0.0' %}
     {% if ip_version == 'ipv4' -%}
       -e MON_IP={{ monitor_address }} \
     {% elif ip_version == 'ipv6' -%}
       -e MON_IP=[{{ monitor_address }}] \
     {% endif -%}
   {% elif ip_version == 'ipv4' -%}
   -e MON_IP={{ hostvars[inventory_hostname]['ansible_' + monitor_interface][ip_version]['address'] }} \
   {% elif ip_version =='ipv6' -%}
   -e MON_IP=[{{ hostvars[inventory_hostname]['ansible_' + monitor_interface][ip_version][0]['address'] }}] \
   {% endif -%}
   {{ ceph_mon_docker_extra_env }} \
   {{ceph_docker_registry }}/{{ ceph_docker_image }}:{{ ceph_docker_image_tag }}
ExecStopPost=-/usr/bin/docker stop ceph-mon-%i
Restart=always
RestartSec=10s
TimeoutStartSec=120
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
