# roles/ceph-mon/templates/ceph-mon.service.j2
[Unit]
Description=Ceph Monitor
After=docker.service

[Service]
EnvironmentFile=-/etc/environment
ExecStartPre=-/usr/bin/{{ container_binary }} rm ceph-mon-%i
ExecStartPre=/bin/sh -c '"$(command -v mkdir)" -p /etc/ceph /var/lib/ceph/mon'
ExecStart=/usr/bin/{{ container_binary }} run --rm --name ceph-mon-%i \
  --memory={{ ceph_mon_docker_memory_limit }} \
{% if (container_binary == 'docker' and ceph_docker_version.split('.')[0] is version_compare('13', '>=')) or container_binary == 'podman' -%}
  --cpus={{ ceph_mon_docker_cpu_limit }} \
{% else -%}
  --cpu-quota={{ ceph_mon_docker_cpu_limit * 100000 }} \
{% endif -%}
  -v /var/lib/ceph:/var/lib/ceph:z \
  -v /etc/ceph:/etc/ceph:z \
  -v /var/run/ceph:/var/run/ceph:z \
  -v /etc/localtime:/etc/localtime:ro \
{% if mon_docker_privileged -%}
  --privileged \
{% endif -%}
{% if mon_docker_net_host -%}
  --net=host \
{% endif -%}
  -e IP_VERSION={{ ip_version[-1:] }} \
  -e MON_IP={{ _current_monitor_address }} \
  -e CLUSTER={{ cluster }} \
  -e FSID={{ fsid }} \
  -e MON_PORT={{ ceph_mon_container_listen_port }} \
  -e CEPH_PUBLIC_NETWORK={{ public_network | regex_replace(' ', '') }} \
  -e CEPH_DAEMON=MON \
  {{ ceph_mon_docker_extra_env }} \
  {{ceph_docker_registry }}/{{ ceph_docker_image }}:{{ ceph_docker_image_tag }}
ExecStop=-/usr/bin/{{ container_binary }} stop ceph-mon-%i
ExecStopPost=-/bin/rm -f /var/run/ceph/{{ cluster }}-mon.{{ monitor_name }}.asok
Restart=always
RestartSec=10s
TimeoutStartSec=120
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
